# services/status-led.service
[Unit]
Description=LED Blinking Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/share/venv/bin/python3 /usr/local/bin/status-led.py
Restart=on-failure

[Install]
WantedBy=multi-user.target

